# etc/supervisor.zway_mqtt_http_bridge.conf
[program:zway_mqtt_http_bridge]
process_name = zway_mqtt_http_bridge
command = python -u /opt/z-way-server/automation/modules/MQTT/zway_mqtt_http_bridge.py